# pi-frame/froggie-frame.service
[Unit]
Description=Froggie Frame Photo Display
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Environment=HOME=%h
Environment=PYTHONUNBUFFERED=1
WorkingDirectory=%h/froggie-frame
ExecStart=/usr/bin/python3 %h/froggie-frame/froggie-frame.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
